# monitor_test.service
[Unit]
Description=Monitor Test Process

[Service]
ExecStart=/root/monitor_test.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target
